# deploy/systemd/mxcwpp-manager.service
[Unit]
Description=Matrix Cloud Security Platform Manager (HTTP API Server)
Documentation=https://github.com/matrixplusio/mxcwpp
After=network.target mysql.service
Wants=mysql.service

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/mxcwpp
ExecStart=/usr/bin/mxcwpp-manager -config /etc/mxcwpp/server.yaml
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=mxcwpp-manager

# 安全设置
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/mxcwpp /var/log/mxcwpp
ReadOnlyPaths=/etc/mxcwpp/server.yaml

# 资源限制
LimitNOFILE=65536
LimitNPROC=4096
MemoryLimit=2G
CPUQuota=200%

# 环境变量
Environment="TZ=Asia/Shanghai"

[Install]
WantedBy=multi-user.target
