# deploy/systemd/timers/mxcwpp-clickhouse-backup.timer
[Unit]
Description=mxcwpp ClickHouse backup timer
Documentation=https://github.com/matrixplusio/mxcwpp

[Timer]
OnCalendar=*-*-* 03:00:00
RandomizedDelaySec=600
Persistent=true
Unit=mxcwpp-clickhouse-backup.service

[Install]
WantedBy=timers.target
